# supervisor/conf.d/caddy.conf
[program:caddy]
command=caddy
autostart=true
autorestart=true
startsecs=10
user = node
environment = HOME="/home/node/", USER="node"
