# supervisor/conf.d/caddy.conf
[program:caddy]
command=caddy
autostart=true
autorestart=true
startsecs=10
user = rtorrent
environment = HOME="/home/rtorrent/", USER="rtorrent"
