# supervisor/conf.d/njr.conf
[program:njr]
user = rtorrent
command = nodejs /home/rtorrent/nodejs-rtorrent/index.js
directory = /home/rtorrent/nodejs-rtorrent
autostart = true
autorestart = true
startsecs=8
